# files/etc/uci-defaults/99-set-password
#!/bin/sh
(echo "14285700"; echo "14285700") | passwd root
exit 0
